# board/FT/TXT/rootfs/etc/init.d/rc.local
#!/bin/sh

#-- start ROBOPro app
su ROBOPro -c 'screen -A -m -d -S ROBOPRO /opt/knobloch/run.sh'

